# init/gourdd.socket
[Unit]
Description=gourdd socket

[Socket]
ListenStream=/run/gourd/gourdd.sock
NoDelay=true

[Install]
WantedBy=sockets.target