# install/gourdd.service
[Unit]
Description=gourdd service
Requires=gourdd.socket

[Service]
ExecStart=/usr/local/bin/gourdd
NonBlocking=true

[Install]
WantedBy=multi-user.target